# ec2-ami-updater/ami-monitor.service
[Unit]
Description=AMI Update Monitor for DDEV Instance
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/ec2-ami-updater
ExecStart=/usr/bin/python3 /home/ubuntu/ec2-ami-updater/ami-monitor.py --instance-id i-1234567890abcdef0 --region us-east-1
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=AWS_REGION=us-east-1
Environment=OLD_INSTANCE_ID=i-1234567890abcdef0

[Install]
WantedBy=multi-user.target
